# infra/ansible/roles/systemd_services/templates/diy-bacnet-server.service.j2
[Unit]
Description=DIY BACnet Server (BACnet + JSON-RPC)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ openfdd_user }}
Group={{ openfdd_group }}
WorkingDirectory={{ diy_bacnet_root }}
EnvironmentFile={{ bench_env_path }}
ExecStart={{ diy_bacnet_root }}/.venv/bin/python -m bacpypes_server.main --public --name {{ diy_bacnet_name }} --instance {{ diy_bacnet_instance }}
Restart=on-failure
RestartSec=3
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
ReadWritePaths={{ diy_bacnet_root }}

[Install]
WantedBy=multi-user.target
