# infra/ansible/roles/systemd_services/templates/easyaso-supervisor.service.j2
[Unit]
Description=easy-aso Supervisor (JSON-RPC)
After=network-online.target diy-bacnet-server.service
Wants=network-online.target

[Service]
Type=simple
User={{ openfdd_user }}
Group={{ openfdd_group }}
WorkingDirectory={{ easyaso_root }}
EnvironmentFile={{ bench_env_path }}
ExecStart={{ easyaso_root }}/.venv/bin/easy-aso-supervisor --host {{ easyaso_supervisor_host }} --port {{ easyaso_supervisor_port }}
Restart=on-failure
RestartSec=3
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
ReadWritePaths={{ easyaso_state_dir }}

[Install]
WantedBy=multi-user.target
