# infra/ansible/roles/systemd_services/templates/openfdd-ui-vite.service.j2
[Unit]
Description=Open-FDD UI (Vite dev server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ openfdd_user }}
Group={{ openfdd_group }}
WorkingDirectory={{ openfdd_root }}/apps/desktop-ui
EnvironmentFile={{ bench_env_path }}
ExecStart={{ npm_path }} run dev -- --host 127.0.0.1 --port {{ openfdd_ui_port }}
Restart=on-failure
RestartSec=3
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
ReadWritePaths={{ openfdd_root }}/apps/desktop-ui

[Install]
WantedBy=multi-user.target
