# openclaw/bench/fake_bacnet_devices/templates/fake_bacnet.service.j2
[Unit]
Description=Open FDD fake BACnet device ({{ device_script }})
After=network.target

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory={{ deploy_path }}
ExecStart={{ deploy_path }}/venv/bin/python {{ deploy_path }}/{{ device_script }} --name {{ device_name }} --instance {{ device_instance }}
Restart=always
RestartSec=10
# No log persistence: stdout/stderr discarded (not written to journal or disk)
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
