# scripts/linux-lan/systemd/openfdd-ui-vite.service
[Unit]
Description=Open-FDD UI (Vite dev server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openfdd
Group=openfdd
WorkingDirectory=/opt/open-fdd/apps/desktop-ui
EnvironmentFile=/etc/openfdd/bench.env
ExecStart=/usr/bin/npm run dev -- --host 127.0.0.1 --port 5173
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
