# bacnet_toolshed/systemd/openfdd-bacnet-poll.service.example
# Example: poll enabled BACnet points into workspace/bacnet/polls/samples.csv
# Install: copy to /etc/systemd/system/openfdd-bacnet-poll.service and edit paths/NIC args.
#
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now openfdd-bacnet-poll.service

[Unit]
Description=Open-FDD BACnet poll driver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openfdd
WorkingDirectory=/opt/open-fdd
Environment=OPENFDD_REPO_ROOT=/opt/open-fdd
Environment=PYTHONPATH=/opt/open-fdd
ExecStart=/opt/open-fdd/.venv/bin/python -m bacnet_toolshed.poll_driver \
  --config /opt/open-fdd/workspace/bacnet/commissioning/points.csv \
  --interval 60 \
  --name OpenFDD --instance 599999 --address 192.168.1.50/24:47808
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
