# bacnet_toolshed/systemd/openfdd-bacnet-commission.service.example
# Example: HTTP agent for BACnet discover jobs (port 8767).
# Copy commission.env.example → workspace/bacnet/commissioning/commission.env first.

[Unit]
Description=Open-FDD BACnet commission agent
After=network-online.target

[Service]
Type=simple
User=openfdd
WorkingDirectory=/opt/open-fdd
Environment=OPENFDD_REPO_ROOT=/opt/open-fdd
Environment=PYTHONPATH=/opt/open-fdd
Environment=OPENFDD_BACNET_COMMISSION_BIND=127.0.0.1
Environment=OPENFDD_BACNET_COMMISSION_PORT=8767
ExecStart=/opt/open-fdd/.venv/bin/python -m bacnet_toolshed.commission_agent
Restart=on-failure

[Install]
WantedBy=multi-user.target
