# infra/ansible/templates/ollama.service.j2
[Unit]
Description=Ollama Local LLM
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=5
# Ollama has no built-in auth — never bind 0.0.0.0 on OT LAN (Tenable CRITICAL).
Environment=OLLAMA_HOST=127.0.0.1:11434
Environment=HOME=/usr/share/ollama

[Install]
WantedBy=multi-user.target
