# infra/ansible/templates/ollama.service.j2
[Unit]
Description=Ollama Local LLM
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=5
Environment=HOME=/usr/share/ollama

[Install]
WantedBy=multi-user.target
