# infra/ansible/templates/openfdd-feather-retention.service.j2
[Unit]
Description=Open-FDD feather store retention + GiB cap enforcement
After=network-online.target

[Service]
Type=oneshot
User={{ ansible_user }}
WorkingDirectory={{ openfdd_remote_dir }}/workspace/api
Environment=OPENFDD_REPO_ROOT={{ openfdd_remote_dir }}
Environment=OPENFDD_WORKSPACE_DIR={{ openfdd_remote_dir }}/workspace
Environment=OFDD_DESKTOP_DATA_DIR={{ openfdd_remote_dir }}/workspace/data
Environment=OFDD_FEATHER_RETENTION_DAYS={{ feather_retention_days | default(90) }}
Environment=OFDD_FEATHER_MAX_GIB={{ feather_max_gib | default(5) }}
Environment=OFDD_FEATHER_TRIM_CHUNK_HOURS={{ feather_trim_chunk_hours | default(24) }}
Environment=PYTHONPATH={{ openfdd_remote_dir }}
ExecStart={{ openfdd_remote_dir }}/.venv/bin/python -m openfdd_bridge.feather_store --maintain
