# infra/ansible/templates/openfdd-building-agent.timer.j2
[Unit]
Description=Run Open-FDD building agent check-in every {{ building_agent_interval_hours | default(6) }}h

[Timer]
OnBootSec=10min
OnUnitActiveSec={{ building_agent_interval_hours | default(6) }}h
Persistent=true
Unit=openfdd-building-agent.service

[Install]
WantedBy=timers.target
