# infra/ansible/templates/openfdd-bacnet-poll.service.j2
[Unit]
Description=Open-FDD BACnet poll driver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory={{ openfdd_remote_dir }}
Environment=OPENFDD_REPO_ROOT={{ openfdd_remote_dir }}
Environment=PYTHONPATH={{ openfdd_remote_dir }}
ExecStart={{ openfdd_remote_dir }}/.venv/bin/python -m bacnet_toolshed.poll_driver \
  --config {{ points_csv_path }} \
  --interval {{ bacnet_poll_interval }} \
  --name {{ bacnet_device_name }} --instance {{ bacnet_instance_id }} \
  --address {{ bacnet_bind_address }} \
{% if bacnet_route_aware | default(false) | bool %}
  --route-aware --network {{ bacnet_network }} --router-ip {{ bacnet_router_ip }} --mstp-net {{ bacnet_mstp_net }} \
{% endif %}
  --output {{ polls_csv_path }}
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
