# infra/ansible/templates/openfdd-bridge.service.j2
[Unit]
Description=Open-FDD operator bridge (FastAPI + Rule Lab)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory={{ openfdd_remote_dir }}/workspace/api
Environment=OPENFDD_REPO_ROOT={{ openfdd_remote_dir }}
Environment=OPENFDD_WORKSPACE_DIR={{ openfdd_remote_dir }}/workspace
Environment=OFDD_DESKTOP_DATA_DIR={{ openfdd_remote_dir }}/workspace/data
Environment=PYTHONPATH={{ openfdd_remote_dir }}/workspace/api:{{ openfdd_remote_dir }}/workspace:{{ openfdd_remote_dir }}
Environment=OPENFDD_BACNET_COMMISSION_URL=http://127.0.0.1:{{ bacnet_commission_port }}
EnvironmentFile=-/etc/openfdd/bridge.secret.env
EnvironmentFile=-{{ openfdd_remote_dir }}/workspace/auth.env.local
EnvironmentFile=-{{ openfdd_remote_dir }}/workspace/mcp.env.local
EnvironmentFile=-{{ openfdd_remote_dir }}/workspace/ollama.env.local
Environment=OFDD_BRIDGE_HOST={{ openfdd_bridge_host }}
Environment=OFDD_BRIDGE_PORT={{ openfdd_bridge_port }}
{% if openfdd_bridge_host in ['0.0.0.0', '::'] %}
Environment=OFDD_CORS_ALLOW_PRIVATE_LAN=1
{% endif %}
ExecStart={{ openfdd_remote_dir }}/.venv/bin/uvicorn openfdd_bridge.main:app --host {{ openfdd_bridge_host }} --port {{ openfdd_bridge_port }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
