# infra/ansible/templates/openfdd-fdd-loop.timer.j2
[Unit]
Description=Run Open-FDD batch FDD every {{ fdd_loop_interval_hours | default(3) }}h

[Timer]
OnBootSec=5min
OnUnitActiveSec={{ fdd_loop_interval_hours | default(3) }}h
Persistent=true
Unit=openfdd-fdd-loop.service

[Install]
WantedBy=timers.target
